# systemd/vnode.target
[Unit]
Description=vNode Multi-Instance Target
After=docker.service
Requires=docker.service

[Install]
WantedBy=multi-user.target
